# packages/backend/com.quinn.scrivener-watcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.quinn.scrivener-watcher</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/node</string>
        <string>/Volumes/Kevin SSD/kev-dev/quinn-writing-studio/node_modules/tsx/dist/cli.mjs</string>
        <string>/Volumes/Kevin SSD/kev-dev/quinn-writing-studio/packages/backend/src/tools/scrivener-watcher.ts</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Volumes/Kevin SSD/kev-dev/quinn-writing-studio/packages/backend</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/kevinvandever/Library/Logs/quinn-scrivener-watcher.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/kevinvandever/Library/Logs/quinn-scrivener-watcher.log</string>

    <key>ThrottleInterval</key>
    <integer>30</integer>
</dict>
</plist>
